# my-website.service
[Unit]
Description=Containers for my website
After=docker.service

[Service]
Type=simple
WorkingDirectory=/home/deployer/actions-runner/_work/my-website/my-website
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always

[Install]
WantedBy=default.target